# packages/odoopbx/odoopbx-0.18.tar.gz/odoopbx-0.18/salt/roots/odoo/templates/odoo.service
[Unit]
Description=Odoo Service
After=network.target

[Service]
User={{ user }}
Type=simple
WorkingDirectory={{ datadir }}
ExecStart={{ venv }}/bin/python3 {{ srcdir }}/odoo-bin --config={{ confdir }}/odoo.conf
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
